# roles/recyclarr/templates/recyclarr.timer.j2
[Unit]
Description=Recyclarr daily TRaSH Guides sync

[Timer]
OnBootSec=15min
OnCalendar=daily
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
